# etc/systemd/system/tgapi.service
[Unit]
Description=TGAPI自助申请
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/tg-api
ExecStart=/usr/bin/python3 /root/tg-api/start_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target